# power_hub_exporter.service
[Unit]
Description=Tank Exporter
After=multi-user.target

[Service]
User=pi
ExecStart=/usr/bin/python /usr/local/bin/power_hub_exporter.py

[Install]
WantedBy=multi-user.target
